# shelly-admin/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Shelly-Admin
# Runs the Nginx daemon
# ==============================================================================

# Wait for Shelly-Admin to become available
#bashio::log.info "NginX is waiting for Shelly-Admin to start..."
#bashio::net.wait_for 43812 localhost 300
#
#bashio::log.info "Starting NginX..."
#exec nginx
